# container-root/etc/profile.d/xdg_runtime_dir.zsh
export XDG_RUNTIME_DIR="$HOME/.run/containers"